# os/ubuntu/schatd2.conf
start on (filesystem and networking) or runlevel [2345]
stop on runlevel [016]

env UID=schatd2
env GID=schatd2
env UMASK=007

respawn
expect fork

exec start-stop-daemon --start --chuid $UID:$GID --umask $UMASK --exec /usr/sbin/schatd2 -- -D