# trunk/os/gentoo/init.d/schatd
#!/sbin/runscript
# $Id$
# Simple Chat
# Copyright © 2008 IMPOMEZIA (http://impomezia.net.ru)

depend() {
	use net
}

start() {
	ebegin "Starting Simple Chat Daemon"
	start-stop-daemon --start --quiet --background --exec /opt/schat/schatd #&>/dev/null
	eend $?
}

stop() {
	ebegin "Stopping Simple Chat Daemon"
	start-stop-daemon --stop --pidfile /opt/schat/schatd.pid
	eend $?
}
